# mqtt-gateway.service
[Unit]
Description=MQTT Gateway Python Service
After=network.target

[Service]
Type=simple
User=gateway
WorkingDirectory=/home/gateway/mqtt-gateway
ExecStart=/home/gateway/mqtt-gateway/venv/bin/python /home/gateway/mqtt-gateway/main.py
Environment=PYTHONUNBUFFERED=1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
